# srv/components/uds/files/uds.service
[Unit]
Description=UDS service
After=network.service

[Service]
Type=simple
User=uds
ExecStart='/usr/bin/uds'

[Install]
